# games/source/supervisord.conf
[supervisord]
nodaemon=true

[program:game_server]
command=/entrypoint.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/game_server.err.log
stdout_logfile=/var/log/game_server.out.log
